# ops/live-node-v1/systemd/atomic-kernel-gates-parity.service
[Unit]
Description=Atomic Kernel Source Capability Parity Gate
After=network-online.target atomic-kernel-mcp.service
Wants=network-online.target atomic-kernel-mcp.service

[Service]
Type=oneshot
WorkingDirectory=/opt/atomic-kernel
EnvironmentFile=/etc/atomic-kernel/node.env
ExecStart=/usr/bin/env npm run -s source:capability:parity:gate
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
